# recipes-connectivity/mbed-cloud-client/mbed-cloud-client-service/mbed-cloud-client.path
[Unit]
Description=Arm MBED Cloud Client path monitor
ConditionPathExists=/etc/systemd/system/multi-user.target.wants/mbed-cloud-client.service

[Path]
PathExists=/var/sota/sota.toml

[Install]
WantedBy=multi-user.target
